# configs/xfce/xfce-overlay/etc/skel/.config/autostart/lxpolkit-manjaro.desktop
[Desktop Entry]
Type=Application
Name=LXPolKit
Exec=lxpolkit
NotShowIn=GNOME;KDE;MATE;Razor;
